# quadlet/git-lan-ssh.container
[Unit]
Description=LAN git server SSH container
Requires=__SERVICE_NAME__.network
After=__SERVICE_NAME__.network

[Quadlet]
DefaultDependencies=false

[Container]
ContainerName=__SERVICE_NAME__-ssh
Image=__IMAGE_TAG__
Network=__SERVICE_NAME__.network
NetworkAlias=__SERVICE_NAME__-ssh
Pull=never
NoNewPrivileges=true
Tmpfs=/run
Tmpfs=/tmp
Volume=__DATA_ROOT__/repos:/var/lib/gitolite/repositories
Volume=__DATA_ROOT__/gitolite-admin:/var/lib/gitolite/.gitolite
Volume=__DATA_ROOT__/state/git-ssh:/var/lib/gitolite/.ssh
Volume=__DATA_ROOT__/state:/var/lib/git-server

[Service]
Restart=on-failure
RestartSec=5
TimeoutStartSec=900

[Install]
WantedBy=default.target